# lk_flow.ini
[program:lk_flow]
process_name = %(program_name)s
command = /root/miniconda3/envs/lk_flow/bin/lk_flow run
directory = /opt/software/lk_flow/
numprocs = 1
autostart = true
autorestart = true
startretries = 3
user = root
redirect_stderr = true
stdout_logfile = /var/log/supervisor/lk_flow/lk_flow.log
stderr_logfile = /var/log/supervisor/lk_flow/lk_flow_err.log
environment = lk_flow_config_path="/etc/lk_flow_config.yaml"
