# tests/supervisor_file.ini
[program:lk_user_http]
process_name = %(program_name)s
command = /opt/virtualenvs/defect_detect/bin/lk_user run
directory = /opt/software/lk_user/
numprocs = 1
autostart = true
autorestart = true
startretries = 3
user = root
redirect_stderr = true
stdout_logfile = /var/log/supervisor/lk_user/lk_user_http.log
stderr_logfile = /var/log/supervisor/lk_user/lk_user_http_err.log
environment = lk_user_config_path="/etc/lk_user_config.yaml"

[program:lk_user_admin]
process_name = %(program_name)s
command = /opt/virtualenvs/defect_detect/bin/lk_user admin
directory = /opt/software/lk_user/
numprocs = 1
autostart = true
autorestart = true
startretries = 3
user = root
redirect_stderr = true
stdout_logfile = /var/log/supervisor/lk_user/lk_user_admin.log
stderr_logfile = /var/log/supervisor/lk_user/lk_user_admin_err.log
environment = lk_user_config_path="/etc/lk_user_config.yaml"

[program:lk_user_grpc]
process_name = %(program_name)s
command = /opt/virtualenvs/defect_detect/bin/lk_user grpc
directory = /opt/software/lk_user/
numprocs = 1
autostart = true
autorestart = true
startretries = 3
user = root
redirect_stderr = true
stdout_logfile = /var/log/supervisor/lk_user/lk_user_grpc.log
stderr_logfile = /var/log/supervisor/lk_user/lk_user_grpc_err.log
environment = lk_user_config_path="/etc/lk_user_config.yaml"
